# crontab
0 0 * * * php /home2/thecent9/www/panel/action_touch_groups.php
